# services/runner-battery.service
[Unit]
Description=Runner battery monitor

[Service]
Type=simple
User=matti
Group=matti
WorkingDirectory=/home/matti/runner_ws
ExecStart=/bin/bash -lc 'source /opt/ros/jazzy/setup.bash && source /home/matti/runner_ws/install/setup.bash && exec ros2 run runner_battery battery_node'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
